# ansible/roles/seafile_reload/files/backup-seafile.service
[Unit]
Description=Backup Seafile Server

[Service]
Type=oneshot
ExecStart=/opt/seafile/backup.sh

[Install]
WantedBy=multi-user.target
